# examples/v4l-server/debian/postinst
#!/bin/sh
set -e

systemctl enable v4l-server.service >/dev/null || true
systemctl start v4l-server.service >/dev/null || true

#DEBHELPER#
