# examples/v4l-server/debian/prerm
#!/bin/sh
set -e

systemctl disable v4l-server.service >/dev/null || true
systemctl stop v4l-server.service >/dev/null || true

#DEBHELPER#
